# roles/reverse_proxy/files/nginx-update-notifier.service
[Unit]
Description=Check for nginx and naxsi updates

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nginx-update-notifier.sh